# resources/docker/app/supervisord.conf
[supervisord]
nodaemon=true

[program:aruna_cache]
user=www-data
directory=/srv/aruna
command=/usr/local/bin/php /srv/aruna/app/run.php cache
redirect_stderr=true
stdout_logfile=syslog
