# resources/docker/app/supervisor_app.conf
[program:aruna_cache]
user=www-data
directory=/srv/aruna
command=/usr/local/bin/php /srv/aruna/app/run.php cache -vvv
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
